# docker-assets/initialize-httpd-auth.service
[Unit]
Description=Initializes Httpd External Authentication
Before=network-pre.target
Wants=network-pre.target
ConditionPathExists=/etc/container-environment
[Service]
Type=oneshot
EnvironmentFile=/etc/container-environment
ExecStart=/usr/bin/initialize-httpd-auth.sh
[Install]
WantedBy=multi-user.target
